# v2-services/syncloud-update.timer
[Unit]
Description=Syncloud update check timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
